# auth-portal/s6-rc.d/nginx/run
#!/bin/sh
# Nginx service

# Run nginx in foreground
exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf
