# auth-portal/s6-rc.d/auth-portal/run
#!/command/execlineb -P

# Auth-Portal Go backend service

# Change to app directory
cd /app

# Set up environment
s6-env AUTH_PORTAL_CONFIG=${AUTH_PORTAL_CONFIG:-/app/configs/auth-portal.yaml}

# Execute auth-portal
exec /app/auth-portal
